# updash@.service
[Unit]
Description=Dashcam backup service

[Service]
Type=simple
ExecStart=/usr/local/bin/updash %I
